# supervisor.conf
[program:scmdb_py]
command=/srv/scmdb_py/1_START.sh
autostart=true
autorestart=true
